# extension/ZenkakuHankaku-Key@www.nyanmo.info/backend/ZenkakuHankakuKeyd.service
[Unit]
Description=ZenkakuHankaku-key extension backend

[Service]
Type=simple
User=root

ExecStart=/server.py
Restart=on-abort

[Install]
WantedBy=multi-user.target

